# adblockrouter.cron
# /etc/cron.d/adblockrouter.cron: crontab entries to auto update adblocking

0 */6 * * *	root    /scripts/adblockrouter/adblockrouter.sh >/dev/null
